# extras/init.d/atmosphere
#!/bin/bash
# Atmosphere init.d
# chkconfig: 345 20 80
# description: Atmosphere as a service


NAME='Atmosphere'
DESC="Deals with the workflow associated with running atmospehre as a service"

case "$1" in
status)
    service apache2 status
    service celeryd status
;;
start)
    printf "Starting $NAME...\n"
    printf "Starting Apache..\n"
    service apache2 start
    printf "Starting Celery..\n"
    service celeryd start
;;
stop)
    printf "Stopping $NAME...\n"
    printf "Stopping Apache..\n"
    service apache2 stop
    printf "Stopping Celery..\n"
    service celeryd stop
;;

restart)
    $0 stop
    $0 start
;;

*)
        echo "Usage: $0 {status|start|stop|restart}"
        exit 1
esac
