# extras/init.d/uwsgi.conf.upstart
# uWSGI - Emperor
#

description     "uWSGI Emperor"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

env LOGTO=/var/log/uwsgi/uwsgi.log
env BINPATH=/opt/env/atmo/bin/uwsgi
env VASSALS=/etc/uwsgi/apps-enabled

exec $BINPATH --emperor $VASSALS --logto $LOGTO